# vuln_images/iptables/too_many_requests.service
[Unit]
Description=Too Many Requests Responser
After=multi-user.target

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/python3 /root/too_many_requests.py

[Install]
WantedBy=multi-user.target