# services/crs/img/files/crs.service
[Unit]
Description=CRS CTF Service

[Service]
User=root
WorkingDirectory=/root
ExecStart=/root/crs
Restart=always

[Install]
WantedBy=multi-user.target
