# services/crs/img/files/rc.local
#!/bin/bash
if [ ! -f /etc/ssh_keys_generated ]; then
    echo "Generating SSH host keys ..."
    ssh-keygen -A
    touch /etc/ssh_keys_generated
fi
